# mission-canvas/systemd/mc-telegram.service
[Unit]
Description=Mission Canvas Telegram Bot
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=mical
Group=mical
WorkingDirectory=./mission-canvas
EnvironmentFile=./mission-canvas/mc_telegram.env
ExecStart=/usr/bin/python3 ./mission-canvas/mc_telegram.py
Restart=always
RestartSec=5
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=full
ReadWritePaths=./mission-canvas/mc_telegram_data
ReadOnlyPaths=.

[Install]
WantedBy=multi-user.target
