# bridges/telegram/systemd/rossi-outbox-deliver.service
[Unit]
Description=Rossi Telegram Relay Outbox Delivery
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=mical
Group=mical
WorkingDirectory=/home/mical/fde
EnvironmentFile=/home/mical/fde/palette/bridges/telegram/systemd/rossi-outbox-deliver.env
ExecStart=/usr/bin/python3 /home/mical/fde/palette/bridges/telegram/outbox_deliver.py --impl-id retail-rossi-store --watch --poll-sec 5
Restart=always
RestartSec=3
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=full
ProtectHome=false
ReadWritePaths=/home/mical/fde/implementations/retail/retail-rossi-store/telegram
ReadOnlyPaths=/home/mical/fde/palette/bridges/telegram
UMask=0077

[Install]
WantedBy=multi-user.target

